# conf/fpm.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm
stdout_logfile=/var/log/nginx/php5-fpm.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
