# conf/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/nginx/nginx.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true